# krack/DEBIAN/test/.config/autostart/skype.desktop
[Desktop Entry]
Type=Application
Exec=skype
Name=skype
Icon=skype
Comment=skype run at start
