# krack/DEBIAN/test/.config/autostart/thunderbird.desktop
[Desktop Entry]
Type=Application
Exec=thunderbird
Name=thunderbird
Icon=thunderbird
Comment=thunderbird run at start
